# salt/glusterfs/files/glusterfs-mount.conf
start on started glusterfs-server
stop on stopping network
stop on starting shutdown

script
  mount /srv
end script

